# packaging/systemd/holly-ux.service
[Unit]
Description=Clyde UX Web Application
After=network.target

[Service]
Type=simple
User=clyde
Group=clyde
WorkingDirectory=/opt/holly-ux
EnvironmentFile=-/etc/default/holly-ux
ExecStart=/opt/holly-ux/.venv/bin/python /opt/holly-ux/main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
